# templates/default/upstart/hello-world.conf.erb
# hello-world - A simple http hello-world service
#
# hello-world replies with a http 200OK and your hostname, simple.
#

respawn
respawn limit 15 5

start on runlevel [2345]
stop on runlevel [06]

kill timeout 60

setuid <%= @user %>

exec <%= @exec_command %>
